# nfpm/deb_files/postinst
#!/bin/sh
set -e

case "$1" in
configure)
	if type systemctl >/dev/null 2>/dev/null; then
		systemctl daemon-reload

		if ! systemctl is-enabled -q hornet.service; then
			systemctl enable hornet.service >/dev/null 2>/dev/null
		fi

		systemctl --no-block restart hornet.service
	fi
	;;

abort-upgrade | abort-remove | abort-deconfigure) ;;

*)
	echo "postinst called with unknown argument \`$1'" >&2
	exit 1
	;;
esac

exit 0
