# nfpm/deb_files/postinst
#!/bin/sh
set -e

case "$1" in
configure)

	chown -R hornet:hornet /var/lib/hornet

	if type systemctl >/dev/null 2>/dev/null; then
		systemctl daemon-reload

		# Only restart if the HORNET service is enabled
		if systemctl is-enabled -q hornet.service; then
			systemctl --no-block restart hornet.service
		fi

	fi
	;;

abort-upgrade | abort-remove | abort-deconfigure) ;;

*)
	echo "postinst called with unknown argument \`$1'" >&2
	exit 1
	;;
esac

exit 0
