# nfpm/rpm_files/prerm
if [ $1 == 0 ]; then #uninstall
    systemctl unmask %{name}.service
    systemctl stop %{name}.service
    systemctl disable %{name}.service
fi

%postun
if [ $1 == 0 ]; then #uninstall
    systemctl daemon-reload
    systemctl reset-failed
fi
