# nfpm/rpm_files/prerm
if [ $1 -eq 0 ]; then
    /usr/bin/systemctl --no-reload disable hornet-legacy.service >/dev/null 2>&1 ||:
    /usr/bin/systemctl stop hornet-legacy.service >/dev/null 2>&1 ||:
fi
